# autostart/Handy.desktop
[Desktop Entry]
Name=Handy
Exec=/usr/bin/handy
Type=Application
Categories=Utility;
X-GNOME-Autostart-enabled=true
